# com.mcgws.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcgws.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/mc-gws/venv/bin/python3</string>
        <string>-m</string>
        <string>mcgws.cli</string>
        <string>weekly</string>
        <string>--notify</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/mc-gws</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday only at 7:28 AM -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>28</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.config/mc-gws/logs/launchd_weekly_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.config/mc-gws/logs/launchd_weekly_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
