# com.mcgws.wrap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcgws.wrap</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/mc-gws/venv/bin/python3</string>
        <string>-m</string>
        <string>mcgws.cli</string>
        <string>wrap</string>
        <string>--notify</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/mc-gws</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Mon-Fri at 4:58 PM (2 min early for 5:00 PM delivery) -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>58</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>58</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>58</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>58</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>58</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.config/mc-gws/logs/launchd_wrap_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.config/mc-gws/logs/launchd_wrap_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
